# dienstplan-update.path
[Unit]
Description=Dienstplan Update Trigger

[Path]
PathExists=/opt/dienstplan/.update-requested

[Install]
WantedBy=multi-user.target
